# asma-bot.service
[Unit]
Description=Robbo Obibot - ASMA Chiptune Radio Bot
After=network.target pulseaudio.service
Wants=pulseaudio.service

[Service]
Type=simple
User=boruta
WorkingDirectory=/home/boruta/robbo-obibot
EnvironmentFile=/home/boruta/robbo-obibot/.env
ExecStart=/home/boruta/robbo-obibot/venv/bin/python3 -u asma-bot.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/boruta/robbo-obibot/bot_output.log
StandardError=append:/home/boruta/robbo-obibot/bot_output.log

[Install]
WantedBy=default.target
