# asma-bot.service
[Unit]
Description=Robbo Obibot - ASMA Chiptune Radio Bot
After=network.target pulseaudio.service
Wants=pulseaudio.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/robbo-obibot-ulimate-chiptune-bot
Environment=DISCORD_BOT_TOKEN=%h/.asma-bot-token
ExecStart=%h/robbo-obibot-ulimate-chiptune-bot/venv/bin/python3 -u asma-bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
