# robbo-obibok-strict.service
[Unit]
Description=Robbo Obibok — Discord Chiptune Radio (strict compatibility mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=boruta
Group=boruta
WorkingDirectory=/home/boruta/robbo-obibot
EnvironmentFile=/home/boruta/robbo-obibot/.env
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/home/boruta/robbo-obibot/venv/bin/python3 -u /home/boruta/robbo-obibot/robbo-obibok-strict.py
ExecStartPre=/usr/bin/env sh -c 'rm -f /home/boruta/robbo-obibot/obibok.pid'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
