# robbo-obibok.service
[Unit]
Description=Robbo Obibok — Discord Chiptune Radio

[Service]
Type=simple
WorkingDirectory=%h/robbo-obibok
EnvironmentFile=%h/robbo-obibok/.env
ExecStart=%h/robbo-obibok/venv/bin/python3 -u %h/robbo-obibok/robbo-obibok.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
