# deploy/robbo-obibok-strict.service
[Unit]
Description=Robbo Obibok — Discord Chiptune Radio (strict compatibility mode)

[Service]
Type=simple
WorkingDirectory=%h/robbo-obibok
EnvironmentFile=%h/robbo-obibok/.env
ExecStart=%h/robbo-obibok/venv/bin/python3 -u %h/robbo-obibok/robbo-obibok-strict.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
